# etc/update-motd.d/91-release-upgrade
#!/bin/sh

exec /usr/lib/update-manager/release-upgrade-motd
